# ship/etc/systemd/system/wlgpio.service
[Service]
ExecStart=gpioset -msignal gpiochip0 2=1
